# ~/Library/LaunchAgents/com.keegan.spy.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.keegan.spy.refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>/absolute/path/Market-Trading-Bot/scripts/refresh_prices.sh >> /absolute/path/Market-Trading-Bot/logs/price_refresh.log 2>&1</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>45</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/absolute/path/Market-Trading-Bot/logs/price_refresh.out</string>
  <key>StandardErrorPath</key><string>/absolute/path/Market-Trading-Bot/logs/price_refresh.err</string>
  <key>RunAtLoad</key><true/>
</dict></plist>
